# ansible/nodejs/supervisord-npm.ini
[program:npm]
directory=/{{ path }}/{{ project }}/{{ app }}
command=npm start
autostart=true
autorestart=true
stderr_logfile={{ pathlog }}/supervisord.npm.err.log
stdout_logfile={{ pathlog }}/supervisord.npm.out.log

[group:gnpm]
programs:npm