# ansible/supervisord-gunicorn.ini
[program:gunicorn]
directory=/{{ path }}/{{ project }}/{{ app }}/
command=/{{ path }}/{{ project }}/{{ app }}//env/bin/gunicorn --workers 3 --bind unix:/{{ path }}/{{ project }}/{{ app }}//trak/app.sock trak.wsgi:application
autostart=true
autorestart=true
stderr_logfile={{ pathlog }}/supervisord.gunicorn.err.log
stdout_logfile={{ pathlog }}/supervisord.gunicorn.out.log

[group:guni]
programs:gunicorn